# docker/etc/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv bash

exec /usr/sbin/cron
