# docker/etc/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv bash

env > /etc/environment

crontab /etc/crontab

exec s6-setuidgid root /usr/sbin/cron -f -L 5
