# docker/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv bash

exec /usr/bin/redis-server /etc/redis/redis.conf --loadmodule /usr/lib/redisearch.so
